# comptetonble/main/conf/crontab/dev
############################
# www.dev-comptetonble.fr #
############################

* * * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble website/cron/certbot:ping
* * * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble website/cron/certbot:rewrite
*/5 * * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble session/cron/sessionClean
0 * * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble mail/cron/mail
0 * * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble selling/cron/invoice
0 * * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble storage/cron/buffer
0 0 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble user/cron/seniority
0 0 1 10 * php /var/www/ctb/framework/lime.php -e dev -a comptetonble company/cron/newSeason
0 12 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble mail/cron/clean
0 2 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble media/cron/clean
0 3 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble hr/cron/calculate
0 4 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble company/cron/demo
0 6 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble cron/sitemap
0 6 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble user/cron/close
0 6 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble user/cron/log
17 7 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble dev/cron/minifyClean
19 5 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble selling/cron/clean
28 4 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble website/cron/certbot:clean
3 */12 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble website/cron/certbot:certificate
47 5 * * * php /var/www/ctb/framework/lime.php -e dev -a comptetonble dev/cron/error