# petitagriculteur/main/conf/crontab/dev
############################
# www.dev-petitagriculteur.fr #
############################

* * * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur website/cron/certbot:ping
* * * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur website/cron/certbot:rewrite
*/5 * * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur session/cron/sessionClean
0 * * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur mail/cron/mail
0 * * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur selling/cron/invoice
0 * * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur storage/cron/buffer
0 0 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur user/cron/seniority
0 0 1 10 * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur company/cron/newSeason
0 12 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur mail/cron/clean
0 2 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur media/cron/clean
0 3 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur hr/cron/calculate
0 4 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur company/cron/demo
0 6 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur cron/sitemap
0 6 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur user/cron/close
0 6 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur user/cron/log
17 7 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur dev/cron/minifyClean
19 5 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur selling/cron/clean
28 4 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur website/cron/certbot:clean
3 */12 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur website/cron/certbot:certificate
47 5 * * * php /var/www/pa/framework/lime.php -e dev -a petitagriculteur dev/cron/error
